# scripts/cz-backend.service
[Unit]
Description=Content Zavod Python backend (FastAPI/uvicorn)
After=network.target

[Service]
Type=simple
User=albert
Group=albert
WorkingDirectory=/home/albert/cz-backend
EnvironmentFile=/home/albert/cz-backend/.env
ExecStart=/home/albert/cz-backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 1 --log-level info
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Безопасность
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/albert/cz-backend
PrivateTmp=true

[Install]
WantedBy=multi-user.target
